# root/etc/s6-overlay/s6-rc.d/init-rtbit-config/run
#!/usr/bin/with-contenv bash

# Ensure directories exist
mkdir -p \
    /home/rtbit/db \
    /home/rtbit/cache \
    /home/rtbit/downloads \
    /home/rtbit/completed

# Set ownership (lsiown handles edge cases on various filesystems)
lsiown abc:abc /home/rtbit
lsiown -R abc:abc /home/rtbit/db /home/rtbit/cache
lsiown abc:abc /home/rtbit/downloads /home/rtbit/completed
